# app/ThirdParty/live_capture.service
[Unit]
Description=Checkin live capture data
After=multi-user.target

[Service]
Restart=always
RestartSec=5s

Type=simple
ExecStart=/usr/bin/python3 /var/www/html/app/ThirdParty/live_capture.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
